# deploy/unity-brain.service
# Unity brain-server — systemd unit (PA.4.7). Operator deploy reference.
#
# The persistent, always-on stateful backend. It holds the brain weights +
# memory and orchestrates the donated-GPU compute pool. Donors (remote
# browsers running compute.html) + the operator connect to it via nginx (see
# nginx-unity-brain.conf). Binds localhost only — ALL external access flows
# through the nginx vhost, which is where Forgejo auth + the public/admin
# split live.
#
# Install:
#   sudo cp deploy/unity-brain.service /etc/systemd/system/
#   (edit User / WorkingDirectory / paths below to match the box)
#   sudo systemctl daemon-reload && sudo systemctl enable --now unity-brain
#   sudo systemctl status unity-brain   # verify
#   journalctl -u unity-brain -f        # tail logs

[Unit]
Description=Unity brain-server (distributed-compute brain backend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as a non-root service user that owns the repo checkout.
User=unity
WorkingDirectory=/opt/unity-brain
ExecStart=/usr/bin/node server/brain-server.js

# PA.4.2 + PA.4.8 environment:
#   UAL_PROXY_AUTH=1  → admin gated on the proxy-vouched X-UAL-User header
#                       (loopback no longer implies admin behind the proxy).
#   BRAIN_BIND        → localhost only; nginx is the sole external entry.
Environment=UAL_PROXY_AUTH=1
Environment=BRAIN_BIND=127.0.0.1
Environment=NODE_ENV=production
# Headless box: don't auto-launch a local Chrome GPU client — donors bring
# the GPUs. (compute.html donors connect over the network.)
Environment=DREAM_NO_AUTO_GPU=1
# Preserve the in-progress curriculum walk across crash/incidental restarts — a
# multi-week walk must NOT be wiped by a restart. autoClearStaleState opts out
# under this flag. A PA.4.8 milestone RESIZE clears the old-size weights itself
# (server/brain-server/gpu.js) so it still re-walks fresh at the new tier; this
# flag only protects against losing training on an incidental restart.
Environment=DREAM_KEEP_STATE=1

Restart=always
RestartSec=5

# ── Resource caps — SHARED BOX (Forgejo + all the lab's git run here too).
# Keep the brain under ~75% of the host (the lab box is 8c/16t · 32 GB) so it
# can NEVER peg the CPU or OOM the box and take Forgejo down. These are cgroup-
# enforced by systemd, independent of any in-app sizing: if the brain ever
# exceeds MemoryMax it is OOM-killed ALONE (the box + Forgejo survive) and
# systemd restarts it (DREAM_KEEP_STATE=1 resumes the walk). The in-app budget
# (DREAM_BRAIN_BUDGET_MB / the RAM-safe auto-cap) sizes the brain well under
# these so the hard cap is only a backstop. Tune the numbers to the box.
# IMPORTANT: systemd does NOT support inline (end-of-line) comments on
# directive lines — `MemoryMax=24G  # note` parses the value as "24G  # note",
# fails validation, and SILENTLY IGNORES the cap (the limit becomes infinity).
# Keep every comment on its OWN line, never after a directive.
MemoryAccounting=true
CPUAccounting=true
# soft limit — kernel reclaims aggressively above this
MemoryHigh=20G
# HARD cap (~75% of 32 GB) — brain OOM-killed alone if exceeded; box + Forgejo survive
MemoryMax=24G
# <=12 of 16 threads — always leaves cores for Forgejo + OS
CPUQuota=1200%
# V8 heap <=16 GB (off-heap ArrayBuffers are still counted by MemoryMax)
Environment=NODE_OPTIONS=--max-old-space-size=16384

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
